# deploy/forty-two-watts.service
[Unit]
Description=forty-two-watts EMS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fredde
WorkingDirectory=/home/fredde/forty-two-watts-go
ExecStart=/home/fredde/forty-two-watts-go/forty-two-watts -config config.yaml -web web
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/fredde/forty-two-watts-go/forty-two-watts.log
StandardError=inherit
# Limit memory growth (RPi has 4GB).
MemoryMax=512M

[Install]
WantedBy=multi-user.target
